# dist/init/linux-systemd/update-caddy.timer
; systemctl start update-caddy.timer

[Unit]
Description=Updates a core executable

[Timer]
; You will want to adjust these:
OnBootSec=10min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
